# profile.d/cpp.sh
# a fancier ninja status line which looks like
# [(18) 114/294  38% :: 21.278]
export NINJA_STATUS_LINE="[(%r) %f/%t %p :: %es] "
